# systemd/deck-saves-daemon.service
[Unit]
Description=Emulator save sync (pull on idle when the NAS is newer, push when a game exits)
After=graphical-session.target

[Service]
Type=simple
ExecStart=/bin/bash %h/deck-saves-daemon.sh
# a logout/shutdown still gets one last push
ExecStop=/bin/bash %h/deck-saves.sh backup
TimeoutStopSec=180
Restart=always
RestartSec=30
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6

[Install]
WantedBy=default.target
