# systemd/loadout-srm.path
[Unit]
Description=Watch for pending Loadout changes

[Path]
PathExists=%h/.loadout-dirty
Unit=loadout-srm.service

[Install]
WantedBy=default.target
