# systemd/loadout-worker.service
[Unit]
Description=Loadout background copy worker
StartLimitIntervalSec=600
StartLimitBurst=10
[Service]
Type=oneshot
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=6
ExecStart=/usr/bin/python3 %h/loadout-worker.py
Restart=on-failure
RestartSec=30
