# systemd/loadout-srm.service
[Unit]
Description=Sync Loadout changes into Steam (SRM refresh)
[Service]
Type=oneshot
TimeoutStartSec=1800
ExecStart=/bin/bash %h/loadout-sync.sh
