# systemd/loadout-worker.path
[Unit]
Description=Watch for queued Loadout copies
[Path]
PathExists=%h/.loadout-queue.json
Unit=loadout-worker.service
[Install]
WantedBy=default.target
